# python/deploy/fleet/agent-memory-dashboard.service
[Unit]
Description=Agent Memory metrics dashboard (always-on)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/cursor-dreaming-memory
ExecStart=/bin/bash -lc 'set -a; [ -f %h/.agent-memory.env ] && . %h/.agent-memory.env; set +a; exec %h/.local/bin/uv run dream-memory serve --host 0.0.0.0 --port 8787'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
