# python/deploy/fleet/agent-memory-eval.service
[Unit]
Description=Dream eval + Linear memory sync (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/dreaming-memory
# Source the shell-style env file (export KEY="val") then run the sync flow.
ExecStart=/bin/bash -lc 'set -a; [ -f %h/.agent-memory.env ] && . %h/.agent-memory.env; set +a; exec %h/.local/bin/uv run python deploy/oci/prefect_flow.py'
